# cups/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/cups/
drwx--x--x	root/root	etc/cups/certs/
-rw-r--r--	root/root	etc/cups/classes.conf
-rw-r--r--	root/root	etc/cups/client.conf
-rw-r--r--	root/root	etc/cups/cupsd.conf
drwxr-xr-x	root/root	etc/cups/interfaces/
-rw-r--r--	root/root	etc/cups/mime.convs
-rw-r--r--	root/root	etc/cups/mime.types
drwxr-xr-x	root/root	etc/cups/ppd/
-rw-r--r--	root/root	etc/cups/printers.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/cups
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/cancel
-rwxr-xr-x	root/root	usr/bin/cups-config
-rwxr-xr-x	root/root	usr/bin/cupstestppd
lrwxrwxrwx	root/root	usr/bin/disable -> ../sbin/accept
lrwxrwxrwx	root/root	usr/bin/enable -> ../sbin/accept
-rwxr-xr-x	root/root	usr/bin/lp
-rwxr-xr-x	root/root	usr/bin/lpoptions
-rwsr-xr-x	daemon/daemon	usr/bin/lppasswd
-rwxr-xr-x	root/root	usr/bin/lpq
-rwxr-xr-x	root/root	usr/bin/lpr
-rwxr-xr-x	root/root	usr/bin/lprm
-rwxr-xr-x	root/root	usr/bin/lpstat
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/cups/
-rw-r--r--	root/root	usr/include/cups/cups.h
-rw-r--r--	root/root	usr/include/cups/http.h
-rw-r--r--	root/root	usr/include/cups/image.h
-rw-r--r--	root/root	usr/include/cups/ipp.h
-rw-r--r--	root/root	usr/include/cups/language.h
-rw-r--r--	root/root	usr/include/cups/md5.h
-rw-r--r--	root/root	usr/include/cups/ppd.h
-rw-r--r--	root/root	usr/include/cups/raster.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/cups/
drwxr-xr-x	root/root	usr/lib/cups/backend/
lrwxrwxrwx	root/root	usr/lib/cups/backend/http -> ipp
-rwxr-xr-x	root/root	usr/lib/cups/backend/ipp
-rwxr-xr-x	root/root	usr/lib/cups/backend/lpd
-rwxr-xr-x	root/root	usr/lib/cups/backend/parallel
-rwxr-xr-x	root/root	usr/lib/cups/backend/scsi
-rwxr-xr-x	root/root	usr/lib/cups/backend/serial
-rwxr-xr-x	root/root	usr/lib/cups/backend/socket
-rwxr-xr-x	root/root	usr/lib/cups/backend/usb
drwxr-xr-x	root/root	usr/lib/cups/cgi-bin/
-rwxr-xr-x	root/root	usr/lib/cups/cgi-bin/admin.cgi
-rwxr-xr-x	root/root	usr/lib/cups/cgi-bin/classes.cgi
-rwxr-xr-x	root/root	usr/lib/cups/cgi-bin/jobs.cgi
-rwxr-xr-x	root/root	usr/lib/cups/cgi-bin/printers.cgi
drwxr-xr-x	root/root	usr/lib/cups/daemon/
-rwxr-xr-x	root/root	usr/lib/cups/daemon/cups-lpd
-rwxr-xr-x	root/root	usr/lib/cups/daemon/cups-polld
drwxr-xr-x	root/root	usr/lib/cups/filter/
-rwxr-xr-x	root/root	usr/lib/cups/filter/gziptoany
-rwxr-xr-x	root/root	usr/lib/cups/filter/hpgltops
-rwxr-xr-x	root/root	usr/lib/cups/filter/imagetops
-rwxr-xr-x	root/root	usr/lib/cups/filter/imagetoraster
-rwxr-xr-x	root/root	usr/lib/cups/filter/pdftops
-rwxr-xr-x	root/root	usr/lib/cups/filter/pstops
-rwxr-xr-x	root/root	usr/lib/cups/filter/rastertodymo
-rwxr-xr-x	root/root	usr/lib/cups/filter/rastertoepson
-rwxr-xr-x	root/root	usr/lib/cups/filter/rastertohp
-rwxr-xr-x	root/root	usr/lib/cups/filter/texttops
-rwxr-xr-x	root/root	usr/lib/libcups.a
lrwxrwxrwx	root/root	usr/lib/libcups.so -> libcups.so.2
-rwxr-xr-x	root/root	usr/lib/libcups.so.2
-rwxr-xr-x	root/root	usr/lib/libcupsimage.a
lrwxrwxrwx	root/root	usr/lib/libcupsimage.so -> libcupsimage.so.2
-rwxr-xr-x	root/root	usr/lib/libcupsimage.so.2
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/backend.1.gz
lrwxrwxrwx	root/root	usr/man/man1/cancel.1.gz -> lp.1.gz
-rw-r--r--	root/root	usr/man/man1/cups-config.1.gz
-rw-r--r--	root/root	usr/man/man1/cupstestppd.1.gz
-rw-r--r--	root/root	usr/man/man1/filter.1.gz
-rw-r--r--	root/root	usr/man/man1/lp.1.gz
-rw-r--r--	root/root	usr/man/man1/lpoptions.1.gz
-rw-r--r--	root/root	usr/man/man1/lppasswd.1.gz
-rw-r--r--	root/root	usr/man/man1/lpq.1.gz
-rw-r--r--	root/root	usr/man/man1/lpr.1.gz
-rw-r--r--	root/root	usr/man/man1/lprm.1.gz
-rw-r--r--	root/root	usr/man/man1/lpstat.1.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/classes.conf.5.gz
-rw-r--r--	root/root	usr/man/man5/cupsd.conf.5.gz
-rw-r--r--	root/root	usr/man/man5/mime.convs.5.gz
-rw-r--r--	root/root	usr/man/man5/mime.types.5.gz
-rw-r--r--	root/root	usr/man/man5/printers.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/accept.8.gz
-rw-r--r--	root/root	usr/man/man8/cups-lpd.8.gz
-rw-r--r--	root/root	usr/man/man8/cups-polld.8.gz
-rw-r--r--	root/root	usr/man/man8/cupsaddsmb.8.gz
-rw-r--r--	root/root	usr/man/man8/cupsd.8.gz
lrwxrwxrwx	root/root	usr/man/man8/disable.8.gz -> enable.8.gz
-rw-r--r--	root/root	usr/man/man8/enable.8.gz
-rw-r--r--	root/root	usr/man/man8/lpadmin.8.gz
-rw-r--r--	root/root	usr/man/man8/lpc.8.gz
-rw-r--r--	root/root	usr/man/man8/lpinfo.8.gz
-rw-r--r--	root/root	usr/man/man8/lpmove.8.gz
lrwxrwxrwx	root/root	usr/man/man8/reject.8.gz -> accept.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/accept
-rwxr-xr-x	root/root	usr/sbin/cupsaddsmb
-rwxr-xr-x	root/root	usr/sbin/cupsd
-rwxr-xr-x	root/root	usr/sbin/lpadmin
-rwxr-xr-x	root/root	usr/sbin/lpc
-rwxr-xr-x	root/root	usr/sbin/lpinfo
-rwxr-xr-x	root/root	usr/sbin/lpmove
lrwxrwxrwx	root/root	usr/sbin/reject -> accept
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/cups/
drwxr-xr-x	root/root	usr/share/cups/banners/
-rw-r--r--	root/root	usr/share/cups/banners/classified
-rw-r--r--	root/root	usr/share/cups/banners/confidential
-rw-r--r--	root/root	usr/share/cups/banners/secret
-rw-r--r--	root/root	usr/share/cups/banners/standard
-rw-r--r--	root/root	usr/share/cups/banners/topsecret
-rw-r--r--	root/root	usr/share/cups/banners/unclassified
drwxr-xr-x	root/root	usr/share/cups/charsets/
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-1
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-10
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-13
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-14
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-15
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-2
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-3
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-4
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-5
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-6
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-7
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-8
-rw-r--r--	root/root	usr/share/cups/charsets/iso-8859-9
-rw-r--r--	root/root	usr/share/cups/charsets/koi8-r
-rw-r--r--	root/root	usr/share/cups/charsets/koi8-u
-rw-r--r--	root/root	usr/share/cups/charsets/utf-8
-rw-r--r--	root/root	usr/share/cups/charsets/windows-1250
-rw-r--r--	root/root	usr/share/cups/charsets/windows-1251
-rw-r--r--	root/root	usr/share/cups/charsets/windows-1252
-rw-r--r--	root/root	usr/share/cups/charsets/windows-1253
-rw-r--r--	root/root	usr/share/cups/charsets/windows-1254
-rw-r--r--	root/root	usr/share/cups/charsets/windows-1255
-rw-r--r--	root/root	usr/share/cups/charsets/windows-1256
-rw-r--r--	root/root	usr/share/cups/charsets/windows-1257
-rw-r--r--	root/root	usr/share/cups/charsets/windows-1258
-rw-r--r--	root/root	usr/share/cups/charsets/windows-874
drwxr-xr-x	root/root	usr/share/cups/data/
-rw-r--r--	root/root	usr/share/cups/data/HPGLprolog
-rw-r--r--	root/root	usr/share/cups/data/psglyphs
-rw-r--r--	root/root	usr/share/cups/data/testprint.ps
drwxr-xr-x	root/root	usr/share/cups/doc/
-rw-r--r--	root/root	usr/share/cups/doc/cmp.html
-rw-r--r--	root/root	usr/share/cups/doc/cups.css
-rw-r--r--	root/root	usr/share/cups/doc/cupsdoc.css
-rw-r--r--	root/root	usr/share/cups/doc/documentation.html
-rw-r--r--	root/root	usr/share/cups/doc/idd.html
drwxr-xr-x	root/root	usr/share/cups/doc/images/
-rw-r--r--	root/root	usr/share/cups/doc/images/accept-jobs.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/add-class.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/add-printer.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/cancel-job.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/classes.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/config-printer.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/continue.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/cups-block-diagram.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/cups-large.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/cups-medium.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/cups-small.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/delete-class.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/delete-printer.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/hold-job.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/left.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/logo.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/manage-classes.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/manage-jobs.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/manage-printers.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/modify-class.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/modify-printer.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/navbar.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/print-test-page.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/printer-idle.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/printer-processing.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/printer-stopped.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/reject-jobs.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/release-job.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/restart-job.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/right.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/set-as-default.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/show-active.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/show-completed.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/start-class.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/start-printer.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/stop-class.gif
-rw-r--r--	root/root	usr/share/cups/doc/images/stop-printer.gif
-rw-r--r--	root/root	usr/share/cups/doc/index.html
-rw-r--r--	root/root	usr/share/cups/doc/ipp.html
-rw-r--r--	root/root	usr/share/cups/doc/overview.html
-rw-r--r--	root/root	usr/share/cups/doc/robots.txt
-rw-r--r--	root/root	usr/share/cups/doc/sam.html
-rw-r--r--	root/root	usr/share/cups/doc/sdd.html
-rw-r--r--	root/root	usr/share/cups/doc/spm.html
-rw-r--r--	root/root	usr/share/cups/doc/sps.html
-rw-r--r--	root/root	usr/share/cups/doc/ssr.html
-rw-r--r--	root/root	usr/share/cups/doc/stp.html
-rw-r--r--	root/root	usr/share/cups/doc/sum.html
-rw-r--r--	root/root	usr/share/cups/doc/svd.html
-rw-r--r--	root/root	usr/share/cups/doc/translation.html
drwxr-xr-x	root/root	usr/share/cups/fonts/
-rw-r--r--	root/root	usr/share/cups/fonts/Courier
-rw-r--r--	root/root	usr/share/cups/fonts/Courier-Bold
-rw-r--r--	root/root	usr/share/cups/fonts/Courier-BoldOblique
-rw-r--r--	root/root	usr/share/cups/fonts/Courier-Oblique
-rw-r--r--	root/root	usr/share/cups/fonts/Symbol
drwxr-xr-x	root/root	usr/share/cups/model/
-rw-r--r--	root/root	usr/share/cups/model/deskjet.ppd
-rw-r--r--	root/root	usr/share/cups/model/deskjet2.ppd
-rw-r--r--	root/root	usr/share/cups/model/dymo.ppd
-rw-r--r--	root/root	usr/share/cups/model/epson24.ppd
-rw-r--r--	root/root	usr/share/cups/model/epson9.ppd
-rw-r--r--	root/root	usr/share/cups/model/laserjet.ppd
-rw-r--r--	root/root	usr/share/cups/model/okidat24.ppd
-rw-r--r--	root/root	usr/share/cups/model/okidata9.ppd
-rw-r--r--	root/root	usr/share/cups/model/stcolor.ppd
-rw-r--r--	root/root	usr/share/cups/model/stcolor2.ppd
-rw-r--r--	root/root	usr/share/cups/model/stphoto.ppd
-rw-r--r--	root/root	usr/share/cups/model/stphoto2.ppd
-rw-r--r--	root/root	usr/share/cups/model/zebra.ppd
drwxr-xr-x	root/root	usr/share/cups/templates/
-rw-r--r--	root/root	usr/share/cups/templates/add-class.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/add-printer.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/admin-op.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/admin.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/choose-device.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/choose-make.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/choose-members.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/choose-model.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/choose-serial.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/choose-uri.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/class-added.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/class-confirm.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/class-deleted.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/class-modified.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/classes.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/config-printer.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/config-printer2.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/error.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/header.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/job-cancel.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/job-hold.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/job-op.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/job-release.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/job-restart.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/jobs.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/modify-class.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/modify-printer.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/option-boolean.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/option-conflict.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/option-header.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/option-pickmany.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/option-pickone.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/option-trailer.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-accept.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-added.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-configured.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-confirm.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-default.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-deleted.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-modified.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-purge.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-reject.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-start.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printer-stop.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/printers.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/test-page.tmpl
-rw-r--r--	root/root	usr/share/cups/templates/trailer.tmpl
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/cups/
drwxr-xr-x	root/root	var/spool/
drwx--x---	root/root	var/spool/cups/
drwxrwx--T	root/root	var/spool/cups/tmp/
