# hal/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/hald
drwxr-xr-x	root/root	etc/udev/
drwxr-xr-x	root/root	etc/udev/rules.d/
-rw-r--r--	root/root	etc/udev/rules.d/90-hal.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/hal-device
-rwxr-xr-x	root/root	usr/bin/hal-find-by-capability
-rwxr-xr-x	root/root	usr/bin/hal-find-by-property
-rwxr-xr-x	root/root	usr/bin/hal-get-property
-rwxr-xr-x	root/root	usr/bin/hal-set-property
-rwxr-xr-x	root/root	usr/bin/lshal
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/dbus-1/
drwxr-xr-x	root/root	usr/etc/dbus-1/system.d/
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/hal.conf
drwxr-xr-x	root/root	usr/etc/hal/
drwxr-xr-x	root/root	usr/etc/hal/fdi/
drwxr-xr-x	root/root	usr/etc/hal/fdi/information/
drwxr-xr-x	root/root	usr/etc/hal/fdi/policy/
drwxr-xr-x	root/root	usr/etc/hal/fdi/preprobe/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/hal/
-rw-r--r--	root/root	usr/include/hal/libhal-storage.h
-rw-r--r--	root/root	usr/include/hal/libhal.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/hal/
-rwxr-xr-x	root/root	usr/lib/hal/hal-system-power-pmu
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-acpi
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-hid-ups
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-keyboard
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-pmu
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-storage
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-usb-csr
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-hiddev
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-input
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-pc-floppy
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-printer
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-serial
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-smbios
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-storage
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-volume
-rwxr-xr-x	root/root	usr/lib/hal/hald-runner
-rw-r--r--	root/root	usr/lib/libhal-storage.a
-rwxr-xr-x	root/root	usr/lib/libhal-storage.la
lrwxrwxrwx	root/root	usr/lib/libhal-storage.so -> libhal-storage.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libhal-storage.so.1 -> libhal-storage.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libhal-storage.so.1.0.0
-rw-r--r--	root/root	usr/lib/libhal.a
-rwxr-xr-x	root/root	usr/lib/libhal.la
lrwxrwxrwx	root/root	usr/lib/libhal.so -> libhal.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libhal.so.1 -> libhal.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libhal.so.1.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/hal-storage.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/hal.pc
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/hald
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/hal/
drwxr-xr-x	root/root	usr/share/hal/fdi/
-rw-r--r--	root/root	usr/share/hal/fdi/fdi.dtd
-rw-r--r--	root/root	usr/share/hal/fdi/fdi.rng
drwxr-xr-x	root/root	usr/share/hal/fdi/information/
drwxr-xr-x	root/root	usr/share/hal/fdi/information/10freedesktop/
-rw-r--r--	root/root	usr/share/hal/fdi/information/10freedesktop/10-camera-ptp.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/information/10freedesktop/10-cd-dvd-burner.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/information/10freedesktop/10-usb-card-readers.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/information/10freedesktop/10-usb-music-players.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/information/10freedesktop/10-usb-pda.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/information/10freedesktop/10-usb-zip-drives.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/information/10freedesktop/10-wireless-mice.fdi
drwxr-xr-x	root/root	usr/share/hal/fdi/information/20thirdparty/
drwxr-xr-x	root/root	usr/share/hal/fdi/policy/
drwxr-xr-x	root/root	usr/share/hal/fdi/policy/10osvendor/
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-keyboard-policy.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-laptop-panel-mgmt-policy.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-power-mgmt-policy.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/15-storage-luks.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/20-storage-methods.fdi
drwxr-xr-x	root/root	usr/share/hal/fdi/policy/20thirdparty/
drwxr-xr-x	root/root	usr/share/hal/fdi/preprobe/
drwxr-xr-x	root/root	usr/share/hal/fdi/preprobe/10osvendor/
-rw-r--r--	root/root	usr/share/hal/fdi/preprobe/10osvendor/10-ide-drives.fdi
drwxr-xr-x	root/root	usr/share/hal/fdi/preprobe/20thirdparty/
drwxr-xr-x	root/root	usr/share/hal/scripts/
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-luks-remove
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-luks-setup
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-luks-teardown
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-lcd-get-brightness
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-lcd-set-brightness
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-power-hibernate
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-power-reboot
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-power-set-power-save
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-power-shutdown
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-power-suspend
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-storage-cleanup-mountpoint
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-storage-cleanup-mountpoints
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-storage-eject
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-storage-mount
-rwxr-xr-x	root/root	usr/share/hal/scripts/hal-system-storage-unmount
