# pcmcia-cs/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/pcmcia/
drwxr-xr-x	root/root	etc/pcmcia/cis/
-rw-r--r--	root/root	etc/pcmcia/cis/3CCFEM556.dat
-rw-r--r--	root/root	etc/pcmcia/cis/3CXEM556.dat
-rw-r--r--	root/root	etc/pcmcia/cis/COMpad2.dat
-rw-r--r--	root/root	etc/pcmcia/cis/COMpad4.dat
-rw-r--r--	root/root	etc/pcmcia/cis/DP83903.dat
-rw-r--r--	root/root	etc/pcmcia/cis/E-CARD.dat
-rw-r--r--	root/root	etc/pcmcia/cis/LA-PCM.dat
-rw-r--r--	root/root	etc/pcmcia/cis/MT5634ZLX.dat
-rw-r--r--	root/root	etc/pcmcia/cis/NE2K.dat
-rw-r--r--	root/root	etc/pcmcia/cis/PCMLM28.dat
-rw-r--r--	root/root	etc/pcmcia/cis/PE-200.dat
-rw-r--r--	root/root	etc/pcmcia/cis/PE520.dat
-rw-r--r--	root/root	etc/pcmcia/cis/RS-COM-2P.dat
-rw-r--r--	root/root	etc/pcmcia/cis/tamarack.dat
-r--r--r--	root/root	etc/pcmcia/config
-r--r--r--	root/root	etc/pcmcia/config.opts
-r-xr-xr-x	root/root	etc/pcmcia/ftl
-r--r--r--	root/root	etc/pcmcia/ftl.opts
-r-xr-xr-x	root/root	etc/pcmcia/ide
-r--r--r--	root/root	etc/pcmcia/ide.opts
-r-xr-xr-x	root/root	etc/pcmcia/ieee1394
-r--r--r--	root/root	etc/pcmcia/ieee1394.opts
-r-xr-xr-x	root/root	etc/pcmcia/memory
-r--r--r--	root/root	etc/pcmcia/memory.opts
-r-xr-xr-x	root/root	etc/pcmcia/network
-r--r--r--	root/root	etc/pcmcia/network.opts
-r-xr-xr-x	root/root	etc/pcmcia/parport
-r--r--r--	root/root	etc/pcmcia/parport.opts
-r-xr-xr-x	root/root	etc/pcmcia/scsi
-r--r--r--	root/root	etc/pcmcia/scsi.opts
-r-xr-xr-x	root/root	etc/pcmcia/serial
-r--r--r--	root/root	etc/pcmcia/serial.opts
-r--r--r--	root/root	etc/pcmcia/shared
-r-xr-xr-x	root/root	etc/pcmcia/wireless
-r--r--r--	root/root	etc/pcmcia/wireless.opts
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/pcmcia
drwxr-xr-x	root/root	sbin/
-rwsr-xr-x	root/root	sbin/cardctl
-rwxr-xr-x	root/root	sbin/cardmgr
-rwxr-xr-x	root/root	sbin/dump_cis
-rwxr-xr-x	root/root	sbin/pack_cis
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man5/
-r--r--r--	root/root	usr/man/man5/pcmcia.5.gz
-r--r--r--	root/root	usr/man/man5/stab.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-r--r--r--	root/root	usr/man/man8/cardctl.8.gz
-r--r--r--	root/root	usr/man/man8/cardmgr.8.gz
-r--r--r--	root/root	usr/man/man8/dump_cis.8.gz
-r--r--r--	root/root	usr/man/man8/pack_cis.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/pcmcia/
