# linux-identd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/identd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/identd.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/identd
lrwxrwxrwx	root/root	usr/sbin/in.identd -> identd
