# webfs/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/webfsd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/webfsd
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/webfsd.1.gz
