# squid/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/squid
drwxr-xr-x	root/root	etc/squid/
-rw-r--r--	root/root	etc/squid/cachemgr.conf
-rw-r--r--	root/root	etc/squid/cachemgr.conf.default
-rw-r--r--	root/root	etc/squid/mime.conf
-rw-r--r--	root/root	etc/squid/mime.conf.default
-rw-r--r--	root/root	etc/squid/squid.conf
-rw-r--r--	root/root	etc/squid/squid.conf.default
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/squidclient
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/squid/
-rwxr-xr-x	root/root	usr/lib/squid/cachemgr.cgi
-rwxr-xr-x	root/root	usr/lib/squid/diskd
-rwxr-xr-x	root/root	usr/lib/squid/unlinkd
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/cachemgr.cgi.8.gz
-rw-r--r--	root/root	usr/man/man8/squid.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/squid
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/squid/
drwxr-xr-x	root/root	usr/share/squid/errors/
-rw-r--r--	root/root	usr/share/squid/errors/ERR_ACCESS_DENIED
-rw-r--r--	root/root	usr/share/squid/errors/ERR_CACHE_ACCESS_DENIED
-rw-r--r--	root/root	usr/share/squid/errors/ERR_CACHE_MGR_ACCESS_DENIED
-rw-r--r--	root/root	usr/share/squid/errors/ERR_CANNOT_FORWARD
-rw-r--r--	root/root	usr/share/squid/errors/ERR_CONNECT_FAIL
-rw-r--r--	root/root	usr/share/squid/errors/ERR_DNS_FAIL
-rw-r--r--	root/root	usr/share/squid/errors/ERR_ESI
-rw-r--r--	root/root	usr/share/squid/errors/ERR_FORWARDING_DENIED
-rw-r--r--	root/root	usr/share/squid/errors/ERR_FTP_DISABLED
-rw-r--r--	root/root	usr/share/squid/errors/ERR_FTP_FAILURE
-rw-r--r--	root/root	usr/share/squid/errors/ERR_FTP_FORBIDDEN
-rw-r--r--	root/root	usr/share/squid/errors/ERR_FTP_NOT_FOUND
-rw-r--r--	root/root	usr/share/squid/errors/ERR_FTP_PUT_CREATED
-rw-r--r--	root/root	usr/share/squid/errors/ERR_FTP_PUT_ERROR
-rw-r--r--	root/root	usr/share/squid/errors/ERR_FTP_PUT_MODIFIED
-rw-r--r--	root/root	usr/share/squid/errors/ERR_FTP_UNAVAILABLE
-rw-r--r--	root/root	usr/share/squid/errors/ERR_ICAP_FAILURE
-rw-r--r--	root/root	usr/share/squid/errors/ERR_INVALID_REQ
-rw-r--r--	root/root	usr/share/squid/errors/ERR_INVALID_RESP
-rw-r--r--	root/root	usr/share/squid/errors/ERR_INVALID_URL
-rw-r--r--	root/root	usr/share/squid/errors/ERR_LIFETIME_EXP
-rw-r--r--	root/root	usr/share/squid/errors/ERR_NO_RELAY
-rw-r--r--	root/root	usr/share/squid/errors/ERR_ONLY_IF_CACHED_MISS
-rw-r--r--	root/root	usr/share/squid/errors/ERR_READ_ERROR
-rw-r--r--	root/root	usr/share/squid/errors/ERR_READ_TIMEOUT
-rw-r--r--	root/root	usr/share/squid/errors/ERR_SHUTTING_DOWN
-rw-r--r--	root/root	usr/share/squid/errors/ERR_SOCKET_FAILURE
-rw-r--r--	root/root	usr/share/squid/errors/ERR_TOO_BIG
-rw-r--r--	root/root	usr/share/squid/errors/ERR_UNSUP_REQ
-rw-r--r--	root/root	usr/share/squid/errors/ERR_URN_RESOLVE
-rw-r--r--	root/root	usr/share/squid/errors/ERR_WRITE_ERROR
-rw-r--r--	root/root	usr/share/squid/errors/ERR_ZERO_SIZE_OBJECT
drwxr-xr-x	root/root	usr/share/squid/icons/
-rw-r--r--	root/root	usr/share/squid/icons/anthony-binhex.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-bomb.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-box.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-box2.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-c.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-compressed.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-dir.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-dirup.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-dvi.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-f.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-image.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-image2.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-layout.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-link.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-movie.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-pdf.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-portal.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-ps.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-quill.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-script.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-sound.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-tar.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-tex.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-text.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-unknown.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-xbm.gif
-rw-r--r--	root/root	usr/share/squid/icons/anthony-xpm.gif
-rw-r--r--	root/root	usr/share/squid/mib.txt
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwx------	squid/squid	var/log/squid/
drwx------	squid/squid	var/squid/
