# ndiswrapper/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/ndiswrapper/
drwxr-xr-x	root/root	etc/rc.d/
-rwx------	root/root	etc/rc.d/wlan
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/modules/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/misc/
-rw-r--r--	root/root	lib/modules/<kernel-version>/misc/ndiswrapper.ko
drwxr-xr-x	root/root	sbin/
-rwxr-xr-x	root/root	sbin/loadndisdriver
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/ndiswrapper.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/ndiswrapper
-rwxr-xr-x	root/root	usr/sbin/ndiswrapper-buginfo
