# shorewall/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-r-xr--r--	root/root	etc/rc.d/shorewall
drwxr-xr-x	root/root	etc/shorewall/
-rw-------	root/root	etc/shorewall/Makefile
-rw-------	root/root	etc/shorewall/accounting
-rw-------	root/root	etc/shorewall/actions
-rw-------	root/root	etc/shorewall/blacklist
-rw-------	root/root	etc/shorewall/continue
-rw-------	root/root	etc/shorewall/ecn
-rw-------	root/root	etc/shorewall/hosts
-rw-------	root/root	etc/shorewall/init
-rw-------	root/root	etc/shorewall/initdone
-rw-------	root/root	etc/shorewall/interfaces
-rw-------	root/root	etc/shorewall/ipsec
-rw-------	root/root	etc/shorewall/maclist
-rw-------	root/root	etc/shorewall/masq
-rw-------	root/root	etc/shorewall/modules
-rw-------	root/root	etc/shorewall/nat
-rw-------	root/root	etc/shorewall/netmap
-rw-------	root/root	etc/shorewall/params
-rw-------	root/root	etc/shorewall/policy
-rw-------	root/root	etc/shorewall/providers
-rw-------	root/root	etc/shorewall/proxyarp
-rw-------	root/root	etc/shorewall/routestopped
-rw-------	root/root	etc/shorewall/rules
-rwxr--r--	root/root	etc/shorewall/shorewall.conf
-rw-------	root/root	etc/shorewall/start
-rw-------	root/root	etc/shorewall/started
-rw-------	root/root	etc/shorewall/stop
-rw-------	root/root	etc/shorewall/stopped
-rw-------	root/root	etc/shorewall/tcclasses
-rw-------	root/root	etc/shorewall/tcdevices
-rw-------	root/root	etc/shorewall/tcrules
-rw-------	root/root	etc/shorewall/tos
-rw-------	root/root	etc/shorewall/tunnels
-rwxr--r--	root/root	etc/shorewall/zones
drwxr-xr-x	root/root	sbin/
-r-xr--r--	root/root	sbin/shorewall
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/shorewall/
-rw-------	root/root	usr/share/shorewall/Limit
-rw-------	root/root	usr/share/shorewall/action.Drop
-rw-------	root/root	usr/share/shorewall/action.Limit
-rw-------	root/root	usr/share/shorewall/action.Reject
-rw-------	root/root	usr/share/shorewall/action.template
-rw-------	root/root	usr/share/shorewall/actions.std
-rw-------	root/root	usr/share/shorewall/configpath
-r-xr--r--	root/root	usr/share/shorewall/firewall
-r--r--r--	root/root	usr/share/shorewall/functions
-r-xr--r--	root/root	usr/share/shorewall/help
-rw-------	root/root	usr/share/shorewall/macro.AllowICMPs
-rw-------	root/root	usr/share/shorewall/macro.Amanda
-rw-------	root/root	usr/share/shorewall/macro.Auth
-rw-------	root/root	usr/share/shorewall/macro.BitTorrent
-rw-------	root/root	usr/share/shorewall/macro.CVS
-rw-------	root/root	usr/share/shorewall/macro.DNS
-rw-------	root/root	usr/share/shorewall/macro.Distcc
-rw-------	root/root	usr/share/shorewall/macro.DropDNSrep
-rw-------	root/root	usr/share/shorewall/macro.DropUPnP
-rw-------	root/root	usr/share/shorewall/macro.Edonkey
-rw-------	root/root	usr/share/shorewall/macro.FTP
-rw-------	root/root	usr/share/shorewall/macro.Gnutella
-rw-------	root/root	usr/share/shorewall/macro.ICQ
-rw-------	root/root	usr/share/shorewall/macro.IMAP
-rw-------	root/root	usr/share/shorewall/macro.LDAP
-rw-------	root/root	usr/share/shorewall/macro.MySQL
-rw-------	root/root	usr/share/shorewall/macro.NNTP
-rw-------	root/root	usr/share/shorewall/macro.NTP
-rw-------	root/root	usr/share/shorewall/macro.NTPbrd
-rw-------	root/root	usr/share/shorewall/macro.PCA
-rw-------	root/root	usr/share/shorewall/macro.POP3
-rw-------	root/root	usr/share/shorewall/macro.Ping
-rw-------	root/root	usr/share/shorewall/macro.PostgreSQL
-rw-------	root/root	usr/share/shorewall/macro.Rdate
-rw-------	root/root	usr/share/shorewall/macro.Rsync
-rw-------	root/root	usr/share/shorewall/macro.SMB
-rw-------	root/root	usr/share/shorewall/macro.SMBswat
-rw-------	root/root	usr/share/shorewall/macro.SMTP
-rw-------	root/root	usr/share/shorewall/macro.SNMP
-rw-------	root/root	usr/share/shorewall/macro.SPAMD
-rw-------	root/root	usr/share/shorewall/macro.SSH
-rw-------	root/root	usr/share/shorewall/macro.SVN
-rw-------	root/root	usr/share/shorewall/macro.Submission
-rw-------	root/root	usr/share/shorewall/macro.Syslog
-rw-------	root/root	usr/share/shorewall/macro.Telnet
-rw-------	root/root	usr/share/shorewall/macro.Trcrt
-rw-------	root/root	usr/share/shorewall/macro.VNC
-rw-------	root/root	usr/share/shorewall/macro.VNCL
-rw-------	root/root	usr/share/shorewall/macro.Web
-rw-------	root/root	usr/share/shorewall/macro.Webmin
-rw-------	root/root	usr/share/shorewall/macro.template
-rw-------	root/root	usr/share/shorewall/rfc1918
-rw-r--r--	root/root	usr/share/shorewall/version
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/shorewall/
