# ypserv/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/ypserv
-rw-r--r--	root/root	etc/ypserv.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/rpcsvc/
-rw-r--r--	root/root	usr/include/rpcsvc/ypxfrd.x
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/yp/
-rwxr-xr-x	root/root	usr/lib/yp/create_printcap
-r-xr-xr-x	root/root	usr/lib/yp/makedbm
-rwxr-xr-x	root/root	usr/lib/yp/match_printcap
-r-xr-xr-x	root/root	usr/lib/yp/mknetid
-rwxr-xr-x	root/root	usr/lib/yp/pwupdate
-r-xr-xr-x	root/root	usr/lib/yp/revnetgroup
-r-xr-xr-x	root/root	usr/lib/yp/yphelper
-rwxr-xr-x	root/root	usr/lib/yp/ypinit
-rwxr-xr-x	root/root	usr/lib/yp/ypxfr
-rwxr-xr-x	root/root	usr/lib/yp/ypxfr_1perday
-rwxr-xr-x	root/root	usr/lib/yp/ypxfr_1perhour
-rwxr-xr-x	root/root	usr/lib/yp/ypxfr_2perday
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/netgroup.5.gz
-rw-r--r--	root/root	usr/man/man5/ypserv.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/makedbm.8.gz
-rw-r--r--	root/root	usr/man/man8/mknetid.8.gz
-rw-r--r--	root/root	usr/man/man8/pwupdate.8.gz
-rw-r--r--	root/root	usr/man/man8/revnetgroup.8.gz
-rw-r--r--	root/root	usr/man/man8/rpc.yppasswdd.8.gz
-rw-r--r--	root/root	usr/man/man8/rpc.ypxfrd.8.gz
-rw-r--r--	root/root	usr/man/man8/yphelper.8.gz
-rw-r--r--	root/root	usr/man/man8/ypinit.8.gz
-rw-r--r--	root/root	usr/man/man8/yppasswdd.8.gz
-rw-r--r--	root/root	usr/man/man8/yppush.8.gz
-rw-r--r--	root/root	usr/man/man8/ypserv.8.gz
-rw-r--r--	root/root	usr/man/man8/ypxfr.8.gz
-rw-r--r--	root/root	usr/man/man8/ypxfrd.8.gz
drwxr-xr-x	root/root	usr/sbin/
-r-xr-xr-x	root/root	usr/sbin/rpc.yppasswdd
-r-xr-xr-x	root/root	usr/sbin/rpc.ypxfrd
-r-xr-xr-x	root/root	usr/sbin/yppush
-r-xr-xr-x	root/root	usr/sbin/ypserv
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/yp/
-rw-r--r--	root/root	var/yp/Makefile
-rw-r--r--	root/root	var/yp/securenets
