# openvpn/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/openvpn/
-rw-r--r--	root/root	etc/openvpn/openvpn.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/openvpn
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/openvpn.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/openvpn
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/empty/
drwxr-xr-x	root/root	var/log/
-rw-r--r--	root/root	var/log/openvpn.log (EMPTY)
drwxr-xr-x	root/root	var/run/
-rw-r--r--	root/root	var/run/openvpn.pid (EMPTY)
