# apache-tomcat/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/tomcat
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/tomcat/
-rw-r--r--	root/root	usr/lib/tomcat/NOTICE
drwxr-xr-x	root/root	usr/lib/tomcat/bin/
-rw-r--r--	root/root	usr/lib/tomcat/bin/bootstrap.jar
-rw-r--r--	root/root	usr/lib/tomcat/bin/catalina-tasks.xml
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/catalina.sh
-rw-r--r--	root/root	usr/lib/tomcat/bin/commons-daemon.jar
-rw-r--r--	root/root	usr/lib/tomcat/bin/commons-logging-api.jar
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/digest.sh
-rw-r--r--	root/root	usr/lib/tomcat/bin/jkstatus-tasks.xml
-rw-r--r--	root/root	usr/lib/tomcat/bin/jmxaccessor-tasks.xml
-rw-r--r--	root/root	usr/lib/tomcat/bin/jsvc.tar.gz
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/setclasspath.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/shutdown.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/startup.sh
-rw-r--r--	root/root	usr/lib/tomcat/bin/tomcat-juli.jar
-rw-r--r--	root/root	usr/lib/tomcat/bin/tomcat-native.tar.gz
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/tool-wrapper.sh
-rw-r--r--	root/root	usr/lib/tomcat/bin/version.sh
drwxr-xr-x	root/root	usr/lib/tomcat/common/
drwxr-xr-x	root/root	usr/lib/tomcat/common/classes/
drwxr-xr-x	root/root	usr/lib/tomcat/common/endorsed/
drwxr-xr-x	root/root	usr/lib/tomcat/common/i18n/
-rw-r--r--	root/root	usr/lib/tomcat/common/i18n/tomcat-i18n-en.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/i18n/tomcat-i18n-es.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/i18n/tomcat-i18n-fr.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/i18n/tomcat-i18n-ja.jar
drwxr-xr-x	root/root	usr/lib/tomcat/common/lib/
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/commons-el.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/jasper-compiler-jdt.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/jasper-compiler.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/jasper-runtime.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/jsp-api.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/naming-factory-dbcp.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/naming-factory.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/naming-resources.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/servlet-api.jar
drwxr-xr-x	root/root	usr/lib/tomcat/conf/
drwxr-xr-x	root/root	usr/lib/tomcat/conf/Catalina/
drwxr-xr-x	root/root	usr/lib/tomcat/conf/Catalina/localhost/
-rw-------	root/root	usr/lib/tomcat/conf/Catalina/localhost/host-manager.xml
-rw-------	root/root	usr/lib/tomcat/conf/Catalina/localhost/manager.xml
-rw-------	root/root	usr/lib/tomcat/conf/catalina.policy
-rw-------	root/root	usr/lib/tomcat/conf/catalina.properties
-rw-------	root/root	usr/lib/tomcat/conf/context.xml
-rw-------	root/root	usr/lib/tomcat/conf/logging.properties
-rw-------	root/root	usr/lib/tomcat/conf/server-minimal.xml
-rw-------	root/root	usr/lib/tomcat/conf/server.xml
-rw-------	root/root	usr/lib/tomcat/conf/tomcat-users.xml
-rw-------	root/root	usr/lib/tomcat/conf/web.xml
drwxr-xr-x	root/root	usr/lib/tomcat/logs/
drwxr-xr-x	root/root	usr/lib/tomcat/server/
drwxr-xr-x	root/root	usr/lib/tomcat/server/classes/
drwxr-xr-x	root/root	usr/lib/tomcat/server/lib/
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-ant-jmx.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-ant.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-cluster.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-optional.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-storeconfig.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/commons-modeler.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-cgi.renametojar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-default.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-invoker.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-ssi.renametojar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-webdav.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-ajp.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-apr.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-coyote.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-http.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-jkstatus-ant.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-util.jar
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/host-manager/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/host-manager/WEB-INF/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/host-manager/WEB-INF/lib/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/WEB-INF/lib/catalina-host-manager.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/WEB-INF/web.xml
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/host-manager.xml
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/host-manager/images/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/images/add.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/images/asf-logo.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/images/code.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/images/design.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/images/docs.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/images/fix.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/images/tomcat.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/images/update.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/images/void.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/host-manager/manager.xml
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/manager/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/manager/WEB-INF/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/manager/WEB-INF/lib/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/WEB-INF/lib/catalina-manager.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/WEB-INF/lib/commons-fileupload-1.0.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/WEB-INF/web.xml
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/html-manager-howto.html
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/manager/images/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/add.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/asf-logo.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/code.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/design.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/docs.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/fix.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/tomcat.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/update.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/void.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/manager-howto.html
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/manager.xml
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/status.xsd
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/xform.xsl
drwxr-xr-x	root/root	usr/lib/tomcat/shared/
drwxr-xr-x	root/root	usr/lib/tomcat/shared/classes/
drwxr-xr-x	root/root	usr/lib/tomcat/shared/lib/
drwxr-xr-x	root/root	usr/lib/tomcat/temp/
-rw-r--r--	root/root	usr/lib/tomcat/temp/bugzilla37035-safeToDelete.tmp (EMPTY)
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/RELEASE-NOTES.txt
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/lib/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/lib/catalina-root.jar
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/web.xml
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/admin/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/admin/index.html
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/asf-logo-wide.gif
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/favicon.ico
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/index.jsp
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/tomcat-power.gif
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/tomcat.gif
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/META-INF/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/balancer/META-INF/context.xml
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/config/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/config/rules.xml
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/lib/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/lib/catalina-balancer.jar
-rw-r--r--	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/web.xml
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/images/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/webdav/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/webdav/WEB-INF/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/webdav/WEB-INF/web.xml
-rw-r--r--	root/root	usr/lib/tomcat/webapps/webdav/index.html
-rw-r--r--	root/root	usr/lib/tomcat/webapps/webdav/tomcat-power.gif
-rw-r--r--	root/root	usr/lib/tomcat/webapps/webdav/tomcat.gif
drwxr-xr-x	root/root	usr/lib/tomcat/work/
